# systemd/phantom-whatsapp-monitor.service
[Unit]
Description=Phantom WhatsApp Monitor (Ninja mode)
Documentation=https://github.com/NinjaTech-AI/phantom
After=network.target phantom-whatsapp-gateway.service
Wants=phantom-whatsapp-gateway.service

[Service]
Type=simple
User=root
WorkingDirectory=/workspace
# Wait for install.sh to finish rsyncing the python modules before starting.
ExecStartPre=/bin/bash -c 'until [ -f /workspace/phantom/whatsapp_interface.py ]; do sleep 1; done'
# Invoke as a module from /workspace so `import phantom` resolves /workspace/phantom
# as a namespace package (finding top-level whatsapp_interface.py). Running as a
# script from /workspace/phantom instead puts the script dir first on sys.path,
# which resolves `phantom` to the nested /workspace/phantom/phantom/ regular
# package (browser toolkit) and causes ModuleNotFoundError on whatsapp_interface.
ExecStart=/usr/local/bin/python -m phantom.monitor_whatsapp
Restart=always
RestartSec=30s

Environment=HOME=/root
Environment=DISPLAY=:99
Environment=PYTHONPATH=/workspace
EnvironmentFile=-/workspace/phantom/.env
StandardOutput=append:/workspace/logs/phantom-whatsapp-monitor.log
StandardError=append:/workspace/logs/phantom-whatsapp-monitor.log
SyslogIdentifier=phantom-whatsapp-monitor

[Install]
WantedBy=multi-user.target
